# demoproject/deploy/chartit.supervisor.conf
[program:chartit]
command=uwsgi
        --socket /tmp/chartit.wsgi.socket
        --chmod-socket
        --master
        --processes 1
        --max-requests 1000
        --buffer-size 32768
        --no-orphans
        --vacuum
        --disable-logging
        --home /work/virtualenvs/chartit/
        --module django.core.handlers.wsgi:WSGIHandler()
environment=DJANGO_SETTINGS_MODULE='demoproject.settings'
user=www-data
autorestart=true
stopsignal=QUIT
